# stremio/.config/autostart/com.stremio.service.desktop
[Desktop Entry]
Type=Application
Name=Stremio Service
Exec=/usr/bin/stremio-service
X-GNOME-Autostart-enabled=true
