# system/etc/systemd/system/nvidia-max-perf.service
# install: sudo cp to /etc/systemd/system/ && sudo systemctl enable --now nvidia-max-perf.service
[Unit]
Description=Lock NVIDIA GPU to maximum performance clocks
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/nvidia-smi -pm 1
ExecStart=/usr/bin/nvidia-smi -lgc 2145
ExecStop=/usr/bin/nvidia-smi -rgc
ExecStop=/usr/bin/nvidia-smi -pm 0

[Install]
WantedBy=multi-user.target
